# systemd/obsidian-mcp.service
[Unit]
Description=Obsidian MCP Server
Documentation=https://github.com/YOUR_USERNAME/obsidian-mcp
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/obsidian-mcp

# Environment — point to your vault mirror and secrets
EnvironmentFile=/home/ec2-user/obsidian-mcp/.env

ExecStart=/usr/bin/node /home/ec2-user/obsidian-mcp/dist/index.js

# Restart policy
Restart=on-failure
RestartSec=5s
StartLimitIntervalSec=60
StartLimitBurst=3

# Logging — captured by journald; view with: journalctl -u obsidian-mcp -f
StandardOutput=journal
StandardError=journal
SyslogIdentifier=obsidian-mcp

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/ec2-user/obsidian-mcp/data /home/ec2-user/.cache

[Install]
WantedBy=multi-user.target
